# .config/systemd/user/displays.service
[Unit]
Description=Display layout switcher

[Service]
ExecStart=/home/hinton/session/scripts/displays

[Install]
WantedBy=default.target